# backend/systemd/webwarden-logprune.timer
[Unit]
Description=Daily webwarden blocked-log prune

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
