# backend/systemd/webwarden-logprune.service
[Unit]
Description=webwarden blocked-log retention prune
After=webwarden-nft.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/webwarden log --prune
